# app/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon = true
logfile = /var/log/supervisor/supervisord.log
pidfile = /var/run/supervisord.pid
childlogdir = /var/log/supervisor/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:tweet-stream]
command = python3 /src/ingest/tweet_streamer.py
autostart = true
autorestart = true
stderr_logfile = /var/log/supervisor/tweet_stream-stderr.log
stdout_logfile = /var/log/supervisor/tweet_stream-stdout.log